# supervisord-bitzeny.conf
[program:bitzenyd]
command=/start-bitzenyd.sh
numprocs=1
autostart=true
autorestart=false
